# asf-update.service
[Unit]
Description=Checks ASF licenses

[Service]
Type=oneshot
ExecStart=/path/to/checker.py
